# doc/iodine-server.socket
[Unit]
Description=Iodine socket

[Socket]
ListenDatagram=53

[Install]
WantedBy=sockets.target
